# auto-credential-updater.service
[Unit]
Description=Auto Credential Updater Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/project_proxy/mac_proxy
ExecStart=/usr/bin/python3 /project_proxy/mac_proxy/auto_credential_updater.py start
Restart=always
RestartSec=10
StandardOutput=append:/project_proxy/mac_proxy/logs/auto_updater.log
StandardError=append:/project_proxy/mac_proxy/logs/auto_updater.error.log

[Install]
WantedBy=multi-user.target

