# gost-7890-monitor.service
[Unit]
Description=Gost 7890 Monitor Service (WARP Backend)
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/project_proxy/mac_proxy
ExecStart=/bin/bash /project_proxy/mac_proxy/gost_7890_monitor.sh start
ExecStop=/bin/bash /project_proxy/mac_proxy/gost_7890_monitor.sh stop
PIDFile=/project_proxy/mac_proxy/logs/gost_7890_monitor.pid
Restart=always
RestartSec=10
StandardOutput=append:/project_proxy/mac_proxy/logs/gost_7890_monitor_systemd.log
StandardError=append:/project_proxy/mac_proxy/logs/gost_7890_monitor_systemd.error.log

[Install]
WantedBy=multi-user.target

