# services/haproxy_7890/com.macproxy.haproxy7890.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macproxy.haproxy7890</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/Ssd/Projects/mac_proxy/services/haproxy_7890/start_haproxy_7890_delayed.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>300</integer>
    
    <key>LaunchOnlyOnce</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Volumes/Ssd/Projects/mac_proxy/services/haproxy_7890/logs/haproxy_7890_launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/Ssd/Projects/mac_proxy/services/haproxy_7890/logs/haproxy_7890_launchd.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/Ssd/Projects/mac_proxy/services/haproxy_7890</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>

