# services/haproxy_7890/haproxy-7890.service
[Unit]
Description=HAProxy 7890 Service (WARP Backend)
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/project_proxy/mac_proxy/services/haproxy_7890
ExecStart=/bin/bash /project_proxy/mac_proxy/services/haproxy_7890/start_haproxy_7890.sh
ExecStop=/bin/bash /project_proxy/mac_proxy/services/haproxy_7890/stop_haproxy_7890.sh
PIDFile=/project_proxy/mac_proxy/services/haproxy_7890/logs/haproxy_7890.pid
Restart=always
RestartSec=10
StandardOutput=append:/project_proxy/mac_proxy/services/haproxy_7890/logs/haproxy_7890_systemd.log
StandardError=append:/project_proxy/mac_proxy/services/haproxy_7890/logs/haproxy_7890_systemd.error.log

[Install]
WantedBy=multi-user.target

